# rootfs/build/services/ss.conf
[program:ss]

command = "/build/bin/start_ss"
autostart = true
autorestart = true
