# rootfs/build/services/v2ray.conf
[program:v2ray]

command = "/build/bin/start_v2ray"
autostart = true
autorestart = true
