# rootfs/build/services/wireguard.conf
[program:wireguard]
command = /build/bin/start_wg
startsecs = 0
autorestart = false
startretries = 1